# plugins/forge/deploy/forge.service
[Unit]
Description=Roxabi Forge local gallery server
After=network.target

[Service]
Type=simple
Environment=FORGE_DIR=%h/.roxabi/forge
WorkingDirectory=%h/.roxabi/forge
ExecStart=%h/.roxabi/forge/run.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target